# activemq5/etc/profile.d/activemq.sh
if [ -e /etc/init.d/activemq ]; then
  /etc/init.d/activemq start
fi